# salt/salt_root/lcd/50-powersave.rules
#SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_ONLINE}=="0", RUN+="/usr/bin/brightness"
SUBSYSTEM=="power_supply", RUN+="/usr/bin/brightness"
